# etc/profile.d/xdg_cache_home.sh
if [ ${LOGNAME} ]; then
  export XDG_CACHE_HOME="/tmp/${LOGNAME}/.cache"
fi

